# packaging/systemd/initrd-unlocked-client@.service.in
# This file is part of unlocked-client

[Unit]
DefaultDependencies=no
Description=Unlocked key provision client in the initrd
Documentation=https://github.com/kalehmann/unlocked-client

[Service]
ExecStartPre=/usr/bin/sleep 5
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/unlocked-client --config @UNLOCKED_CONFIG_DIR@/%i.conf

[X-SystemdTool]
InitrdBinary=/usr/bin/unlocked-client
InitrdPath=/etc/unlocked
