# packaging/systemd/unlocked-client@.socket.in
# This file is part of unlocked-client

[Unit]
Description=Unlocked key provision client socket
Documentation=https://github.com/kalehmann/unlocked-client
PartOf=unlocked-client@.service
Requires=unlocked-client@.service

[Socket]
ListenStream=/run/unlocked-%i.sock

[Install]
WantedBy=sockets.target
