# packaging/systemd/unlocked-client@.socket.in
[Unit]
DefaultDependencies=no

[Socket]
ListenStream=/run/unlocked-%i.sock

[Install]
RequiredBy=cryptsetup-pre.target
