# packaging/systemd/unlocked-client@.service.in
[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/unlocked-client --config @UNLOCKED_CONFIG_DIR@/%i.conf